# conf/blog.socket
[Unit]
Description=blog socket

[Socket]
ListenStream=/home/javier/Projects/blog/blog.sock

[Install]
WantedBy=sockets.target